# OS/workstation/config/systemd/system/zenith-snapper-firstboot.service
[Unit]
Description=ZenithOS Snapper First-Boot Initialization
After=local-fs.target
ConditionPathExists=/etc/zenith-installed
ConditionPathExists=/etc/snapper/configs/root

[Service]
Type=oneshot
ExecStart=/usr/lib/zenith/zenith-snapper-firstboot
RemainAfterExit=yes

[Install]
WantedBy=multi-user.target