# OS/workstation/config/systemd/system/zenith-session-report-watch.service
[Unit]
Description=ZenithOS desktop smoke session report relay
After=display-manager.service sddm.service
Wants=display-manager.service

[Service]
Type=oneshot
ExecStart=/usr/lib/zenith/zenith-session-report-watch
StandardOutput=journal+console
StandardError=journal+console

[Install]
WantedBy=graphical.target
