# OS/workstation/config/xdg/autostart/zenith-greeter.desktop
[Desktop Entry]
Type=Application
Name=Zenith Greeter
Exec=/usr/lib/zenith/zenith-greeter
OnlyShowIn=KDE;
NoDisplay=true
X-KDE-Autostart-enabled=true
